# .config/autostart/kitty.desktop
[Desktop Entry]
Name=Kitty (Desktop 2)
Comment=Terminal
Exec=kitty
Type=Application
Terminal=false
X-KDE-StartupNotify=true
X-KDE-autostart-phase=2
